# config/client_bot.service
[Unit]
Description=StealthNET Client Telegram Bot
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/path/to/your/api
Environment="PATH=/path/to/your/venv/bin"
ExecStart=/path/to/your/venv/bin/python /path/to/your/api/client_bot.py
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

